# mirror.mount
[Mount]
What=/dev/vdb
Where=/mirror
Type=ext4
